# configs/releng/airootfs/etc/systemd/system/post-setup.service
[Unit]
Description=Post setup mechanism

[Service]
Type=oneshot
ExecStart=/usr/local/bin/post-setup

[Install]
WantedBy=multi-user.target